# scripts/ww.service.example
[Unit]
Description=Worldwave AI Agent Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/worldwave
EnvironmentFile=-%h/worldwave/.env
# Prefer venv python + server.py (honors WW_HOST / WW_PORT from .env)
ExecStart=%h/worldwave/.venv/bin/python %h/worldwave/server.py
Restart=always
RestartSec=10
MemoryMax=2G
TimeoutStopSec=30

[Install]
WantedBy=default.target

